# ansible/roles/litellm/templates/litellm.service.j2
# {{ ansible_managed }}
[Unit]
Description=LiteLLM OpenAI-compatible proxy
After=network-online.target
Wants=network-online.target

[Service]
User={{ litellm_user }}
Group={{ litellm_group }}
EnvironmentFile={{ litellm_env_file }}
# HOME under the daemon's own dir, and the same Prisma cache vars the role used at
# provision time so the daemon reads the already-generated client + engine binary
# instead of trying to fetch them. OFFLINE_MODE makes a missing engine a hard error
# rather than a silent (and, under ProtectSystem=strict, failing) download attempt.
Environment=HOME={{ litellm_home }}
Environment=XDG_CACHE_HOME={{ litellm_prisma_cache }}
Environment=PRISMA_BINARY_CACHE_DIR={{ litellm_prisma_cache }}/prisma-python/binaries
Environment=PRISMA_HIDE_UPDATE_MESSAGE=1
Environment=PRISMA_OFFLINE_MODE=true
# Run the venv's console script directly (absolute path — no PATH reliance under
# systemd). Config carries model routing; the port is pinned here. DATABASE_URL,
# LITELLM_MASTER_KEY and STORE_MODEL_IN_DB come from the EnvironmentFile above.
ExecStart={{ litellm_bin }} --config {{ litellm_config_file }} --port {{ litellm_port }}
Restart=on-failure
RestartSec=5
# Prisma migrations already ran at provision time, so startup is just app boot +
# DB connect — but keep generous headroom so a slow first connect isn't treated as
# a failure.
TimeoutStartSec=120

# Hardening: persistent state lives in Postgres. ProtectSystem=strict makes /
# read-only; PrivateTmp gives a scratch /tmp. ReadWritePaths re-opens just the
# daemon's home so Prisma can cache its engine binary (under the venv / $HOME)
# without widening write access anywhere else.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ litellm_home }}

[Install]
WantedBy=multi-user.target
